# oauth.service
[Unit]
Description=OAuth Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=procesor
Group=procesor
WorkingDirectory=/home/procesor/oauth_service

# Updated paths to match your actual virtual environment location
Environment="PYTHONPATH=/home/procesor/oauth_service"
Environment="PATH=/home/procesor/oauth_service/oauth_service/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/procesor/oauth_service/.env

# Updated ExecStart with correct path
ExecStart=/home/procesor/oauth_service/oauth_service/venv/bin/uvicorn oauth_service.main:app --host 0.0.0.0 --port 8000 --workers 4

Restart=always
RestartSec=1
StartLimitInterval=0

# Security measures
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

# Logging
StandardOutput=append:/var/log/oauth_service/oauth.log
StandardError=append:/var/log/oauth_service/error.log

[Install]
WantedBy=multi-user.target